# quota-monitor.service
[Unit]
Description=Mirrowel Proxy Quota Monitor Overlay
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/quota-monitor
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
